# .docker/supervisor/supervisord.conf
[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisord]
nodaemon=true
logfile=/var/www/storage/logs/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/www/storage/logs

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

; ─────────────────────────────────────────────
; Laravel Horizon (manages Redis queues)
; Uncomment after installing Horizon
; ─────────────────────────────────────────────
;[program:horizon]
;process_name=%(program_name)s
;command=php /var/www/artisan horizon
;autostart=true
;autorestart=true
;user=forge
;redirect_stderr=true
;stdout_logfile=/var/www/storage/logs/horizon.log
;stopwaitsecs=3600

; ─────────────────────────────────────────────
; Default queue worker (used before Horizon)
; ─────────────────────────────────────────────
[program:queue-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/artisan queue:work redis --sleep=3 --tries=3 --max-time=3600
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=forge
numprocs=2
redirect_stderr=true
stdout_logfile=/var/www/storage/logs/worker.log
stopwaitsecs=3600

; ─────────────────────────────────────────────
; Laravel Reverb (WebSockets) — optional
; Uncomment when needed
; ─────────────────────────────────────────────
;[program:reverb]
;process_name=%(program_name)s
;command=php /var/www/artisan reverb:start --host=0.0.0.0 --port=8080
;autostart=true
;autorestart=true
;user=forge
;redirect_stderr=true
;stdout_logfile=/var/www/storage/logs/reverb.log
